# etc/rc.local
#!/bin/bash

cpupower frequency-set --governor performance
echo 3072 > /sys/class/rtc/rtc0/max_user_freq
echo 3072 > /proc/sys/dev/hpet/max-user-freq

exit 0